# ChloroFrame/Resources/fullstacksandbox.com.ChloroFrame.Helper.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <!-- Label must match the helper's bundle identifier exactly. -->
    <key>Label</key>
    <string>fullstacksandbox.com.ChloroFrame.Helper</string>

    <!-- SMJobBless installs the binary to /Library/PrivilegedHelperTools/. -->
    <key>Program</key>
    <string>/Library/PrivilegedHelperTools/fullstacksandbox.com.ChloroFrame.Helper</string>

    <!-- On-demand: launchd starts the helper when the Mach service is first looked up. -->
    <key>MachServices</key>
    <dict>
        <key>fullstacksandbox.com.ChloroFrame.Helper</key>
        <true/>
    </dict>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
